# deploy/com.finitebuilds.osint-api.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.finitebuilds.osint-api</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/ves/Developer/OSINT/project/deploy/run-api.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/ves/Developer/OSINT/project</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>StandardOutPath</key>
    <string>/Users/ves/Library/Logs/osint-api.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/ves/Library/Logs/osint-api.log</string>
</dict>
</plist>
